# crontab
51 * * * * python3 /app/nightly_update.py
