# crontab
45 15 * * 1-5 python3 /app/nightly_update.py
